# ippayment/ippayment_upstart.conf
start on (local-filesystems
    and net-device-up IFACE=lo)
    
stop on shutdown

respawn

script

start-stop-daemon --start  --make-pidfile --pidfile /var/run/ippayment.pid -c juno --exec /home/juno/git/ippayment/ippayment/ippayment.bash

end script